# dashboard/deploy/etf-dashboard.service
[Unit]
Description=ETF Rotation Dashboard (uvicorn)
After=network.target

[Service]
Type=simple
User=etf
WorkingDirectory=/opt/etf-dashboard/dashboard/backend
Environment=DATA_DIR=/opt/etf-dashboard/data
ExecStart=/opt/etf-dashboard/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8765 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Uncomment to protect Settings/Recalculate endpoints:
# Environment=DASHBOARD_SECRET=change-me-to-a-strong-secret

[Install]
WantedBy=multi-user.target
